# software/systemd/yamper.service
[Unit]
Description=Yamper AI voice robot
After=network-online.target sound.target
Wants=network-online.target

[Service]
Type=simple
User=pi
Group=pi
WorkingDirectory=/home/pi/yamper/software
EnvironmentFile=/home/pi/yamper/software/.env
Environment=PYTHONUNBUFFERED=1
ExecStart=/home/pi/yamper/software/.venv/bin/python3 -m yamper.main
Restart=on-failure
RestartSec=10
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
